# deploy/sunshine-vd.service
[Unit]
Description=Sunshine Virtual Display Daemon
After=display-manager.service dbus.service
Requires=dbus.service

[Service]
Type=simple
User=root
Group=sunshine-vd
ExecStart=/usr/local/bin/svd-daemon
Restart=on-failure
RestartSec=3
TimeoutStopSec=45

RuntimeDirectory=sunshine-vd
RuntimeDirectoryMode=0750
StateDirectory=sunshine-vd
StateDirectoryMode=0750
UMask=0007

# DRM/debugfs writes plus dropping kscreen-doctor to the desktop UID/GID.
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_SETUID CAP_SETGID
AmbientCapabilities=CAP_SETUID CAP_SETGID
NoNewPrivileges=yes
ProtectSystem=full
PrivateTmp=yes
PrivateNetwork=yes
RestrictAddressFamilies=AF_UNIX
ProtectClock=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
